# net-p2p/go-ipfs/files/go-ipfs.service
[Unit]
Description=IPFS daemon
After=network.target

[Service]
ExecStart=/usr/bin/ipfs daemon
User=go-ipfs
Group=go-ipfs

